# data/lib/systemd/user/ico-uxf-wait-launchpad-ready.path
# Description=Path activation for ico-system-controller

[Path]
PathExists=/tmp/amd_ready
Unit=ico-system-controller.service

[Install]
WantedBy=weston.target
